# systemd/bytecli.service
[Unit]
Description=ByteCLI Voice Dictation Service
After=graphical-session.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=/usr/bin/bytecli-service
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
