# debian/postinst
#!/bin/bash
set -e

# ---------------------------------------------------------------------------
# ByteCLI postinst -- runs as root after dpkg unpacks the package.
# ---------------------------------------------------------------------------
# Responsibilities:
#   1. Install openai-whisper (not available in apt).
#   2. Enable and start the systemd user service.
# ---------------------------------------------------------------------------

case "$1" in
  configure)

    # -----------------------------------------------------------------------
    # 1. Install Whisper via pip (not available in apt).
    # -----------------------------------------------------------------------
    echo "ByteCLI: Installing speech recognition engine (openai-whisper + PyTorch)..."
    echo "ByteCLI: This may take several minutes on first install."
    pip3 install --break-system-packages openai-whisper 2>/dev/null || \
    pip3 install openai-whisper 2>/dev/null || true

    # Install additional Python dependencies not available in apt.
    echo "ByteCLI: Installing Python dependencies (python-xlib, pulsectl, sounddevice)..."
    pip3 install --break-system-packages python-xlib pulsectl sounddevice 2>/dev/null || \
    pip3 install python-xlib pulsectl sounddevice 2>/dev/null || true

    # -----------------------------------------------------------------------
    # 2. Enable the systemd user service.
    # -----------------------------------------------------------------------
    systemctl --global enable bytecli.service 2>/dev/null || true

    # Immediate start for the user who ran sudo.
    if [ -n "$SUDO_USER" ]; then
      SUDO_UID=$(id -u "$SUDO_USER")
      if [ -d "/run/user/$SUDO_UID" ]; then
        su - "$SUDO_USER" -c \
          "systemctl --user daemon-reload && systemctl --user start bytecli.service" \
          2>/dev/null || true
      fi
    fi
    ;;
esac

exit 0
